# ports/sysutils/dbus/pkg/PLIST
@comment $MirOS$
@comment $NetBSD: PLIST,v 1.9 2006/12/17 15:39:33 jmmv Exp $
bin/dbus-cleanup-sockets
bin/dbus-daemon
bin/dbus-launch
bin/dbus-monitor
bin/dbus-send
bin/dbus-uuidgen
include/dbus-1.0/dbus/dbus-address.h
include/dbus-1.0/dbus/dbus-bus.h
include/dbus-1.0/dbus/dbus-connection.h
include/dbus-1.0/dbus/dbus-errors.h
include/dbus-1.0/dbus/dbus-macros.h
include/dbus-1.0/dbus/dbus-memory.h
include/dbus-1.0/dbus/dbus-message.h
include/dbus-1.0/dbus/dbus-misc.h
include/dbus-1.0/dbus/dbus-pending-call.h
include/dbus-1.0/dbus/dbus-protocol.h
include/dbus-1.0/dbus/dbus-server.h
include/dbus-1.0/dbus/dbus-shared.h
include/dbus-1.0/dbus/dbus-signature.h
include/dbus-1.0/dbus/dbus-threads.h
include/dbus-1.0/dbus/dbus-types.h
include/dbus-1.0/dbus/dbus.h
lib/dbus-1.0/include/dbus/dbus-arch-deps.h
lib/libdbus-1.la
lib/pkgconfig/dbus-1.pc
man/man1/dbus-cleanup-sockets.1
man/man1/dbus-daemon.1
man/man1/dbus-launch.1
man/man1/dbus-monitor.1
man/man1/dbus-send.1
man/man1/dbus-uuidgen.1
@sample ${SYSCONFDIR}/dbus-1/
@sample ${SYSCONFDIR}/dbus-1/event.d/
@sample ${SYSCONFDIR}/dbus-1/system.d/
share/examples/dbus/session.conf
@sample ${SYSCONFDIR}/dbus-1/session.conf
share/examples/dbus/system.conf
@sample ${SYSCONFDIR}/dbus-1/system.conf
@comment share/examples/rc.d/messagebus.redhat
@comment share/examples/rc.d/rc.messagebus.slackware
@dirrm share/examples/dbus
@sample var/lib/
@owner _dbus
@group _dbus
@sample var/lib/dbus/
@sample var/run/dbus/
@exec %D/bin/dbus-uuidgen --ensure
@owner
@group
@extraunexec rm -rf %D/var/lib/dbus
@extraunexec rm -rf %D/var/run/dbus
@exec mkdir -p %D/share/dbus-1/services
@dirrm share/dbus-1/services
@dirrm share/dbus-1
@dirrm lib/pkgconfig
@dirrm lib/dbus-1.0/include/dbus
@dirrm lib/dbus-1.0/include
@dirrm lib/dbus-1.0
@dirrm include/dbus-1.0/dbus
@dirrm include/dbus-1.0
